# templates/aries/dlob.hans/mrt3000/udev.rules/60-usb-storage.rules
ACTION=="add", KERNEL=="sd*", SUBSYSTEM=="block", RUN+="/etc/scripts/usbmount.sh add %k"
ACTION=="remove", KERNEL=="sd*", RUN+="/etc/scripts/usbmount.sh remove %k"
